# data/app.conf
start on runlevel [2345] and started docker
stop on runlevel [!2345]

respawn
respawn limit unlimited

post-start script
  stack=$(curl http://169.254.169.254/latest/user-data | jq -r ".signal.stack")
  resource=$(curl http://169.254.169.254/latest/user-data | jq -r ".signal.resource")
  cfn-signal -s true --stack=$stack --resource=$resource
end script

post-stop exec sleep 1

script
  curl http://169.254.169.254/latest/user-data | docker run --sig-proxy -i -v /var/run/docker.sock:/var/run/docker.sock convox/init
end script
